# scripts/rectella-backup.service
[Unit]
Description=Rectella PostgreSQL backup

[Service]
Type=oneshot
ExecStart=/home/bast/brain/projects/code/rectella-shopify-service/scripts/backup.sh
